# packages/openrc/cos-setup/files/etc/init.d/cos-setup-network
#!/sbin/openrc-run

supervise_daemon_args="--stdout /var/log/cos-setup-network.log --stderr /var/log/cos-setup-network.log"

depend() {
  after net
  provide cos-setup-network
}

start() {
  elemental run-stage network
  eend 0
}
