# packages/openrc/cos-setup/files/etc/init.d/cos-setup-network
#!/sbin/openrc-run

depend() {
  after net
  provide cos-setup-network
}

start() {
  elemental run-stage network \
    >>  /var/log/cos-setup-network.log \
    2>> /var/log/cos-setup-network.log
  eend 0
}
