# packages/openrc/cos-setup/files/etc/init.d/cos-setup-boot
#!/sbin/openrc-run

supervise_daemon_args="--stdout /var/log/cos-setup-boot.log --stderr /var/log/cos-setup-network.log"

depend() {
  provide cos-setup-boot
}

start() {
  elemental run-stage boot
  eend 0
}
