# packages/openrc/cos-setup/files/etc/init.d/cos-setup-boot
#!/sbin/openrc-run


depend() {
  provide cos-setup-boot
}

start() {
  elemental \
    --logfile /var/log/cos-setup-boot.log \
    run-stage boot 
  eend 0
}
